# linux_script_monitoramento/monitory.timer
[Unit]
Description=Timer to execute monitory

[Timer]
OnCalendar=*:0/15
Persistent=True

[Install]
WantedBy=timers.target
